# examples/systemd/portal-proxy-prune.timer
[Unit]
Description=Run Portal Proxy pruning daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
